# etc/profile.d/sway-blue.sh
#!/bin/sh
# Alias for sway-blue system recipes
alias sway-just='just --justfile /usr/share/sway-blue/justfile --working-directory ~'

# Default editor
export EDITOR=nvim
export VISUAL=nvim
